# tanlnode@.service
[Unit]
Description=tanlnode

[Service]
Type=simple
Environment=WG_IF=%i
EnvironmentFile=/etc/tanlnode/%i/config
ExecStart=/usr/bin/tanlnode

[Install]
WantedBy=multi-user.target

